# launchd/com.market-monitor.review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-monitor.review</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/xiaosheji/projects/market-monitor/.venv/bin/python</string>
        <string>-m</string>
        <string>market_monitor.cli</string>
        <string>run</string>
        <string>review</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/xiaosheji/projects/market-monitor</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/xiaosheji/projects/market-monitor</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>7</integer>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/com.market-monitor.review.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.market-monitor.review.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
