# launchd/com.market-monitor.hk-shock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-monitor.hk-shock</string>
    <key>ProgramArguments</key>
    <array>
        <string>$PROJECT_ROOT/.venv/bin/python</string>
        <string>-m</string>
        <string>market_monitor.cli</string>
        <string>run</string>
        <string>hk_shock</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$PROJECT_ROOT</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>$PROJECT_ROOT</string>
    </dict>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key><string>/tmp/com.market-monitor.hk-shock.log</string>
    <key>StandardErrorPath</key><string>/tmp/com.market-monitor.hk-shock.err</string>
    <key>RunAtLoad</key><false/>
</dict>
</plist>
